# plant-timelapse.service
[Unit]
Description=Plant Timelapse Monitor
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/plant-timelapse/app.py
WorkingDirectory=/home/pi/plant-timelapse
User=pi
Restart=always
RestartSec=5
# Give the camera stack time to initialize on boot
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
